# systemd/amp-port-sync-oneshot.service
[Unit]
Description=AMP Port Forwarder (oneshot)
Documentation=https://github.com/jtdoepke/amp-mikrotik-port-forwarder
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/amp-port-sync/env
ExecStart=/usr/local/bin/amp-port-sync sync --once

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true

# Allow reading password files
ReadOnlyPaths=/etc/amp-port-sync
